# marketplace/digitalocean/packer/files/etc/update-motd.d/99-defensia
#!/bin/bash
# Defensia Agent — Message of the Day

echo ""
echo "  ╔══════════════════════════════════════════════════╗"
echo "  ║         DEFENSIA Security Agent                  ║"
echo "  ╚══════════════════════════════════════════════════╝"
echo ""

if [[ -f /etc/defensia/config.json ]] && systemctl is-active --quiet defensia-agent; then
    echo "  Status:  ACTIVE ✓"
    echo ""
    echo "  Commands:"
    echo "    systemctl status defensia-agent"
    echo "    journalctl -u defensia-agent -f"
    echo ""
    echo "  Dashboard: https://defensia.cloud/dashboard"
else
    echo "  Status:  NOT CONFIGURED"
    echo ""
    echo "  To activate Defensia on this server:"
    echo ""
    echo "    1. Sign up at https://defensia.cloud"
    echo "    2. Go to Dashboard → Add Server → copy your token"
    echo "    3. Run:"
    echo ""
    echo "       defensia-agent register https://defensia.cloud $(hostname -s) YOUR_TOKEN"
    echo "       systemctl start defensia-agent"
    echo ""
    echo "  Docs: https://github.com/defensia/agent"
fi

echo ""
