# defensia-agent.service
[Unit]
Description=Defensia Security Agent
Documentation=https://defensia.com/docs/agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Auto-restore the binary from backup if it is missing or empty.
# This ensures a failed update never permanently kills the service.
ExecStartPre=/bin/sh -c 'if [ ! -s /usr/local/bin/defensia-agent ] && [ -s /usr/local/bin/defensia-agent.bak ]; then cp /usr/local/bin/defensia-agent.bak /usr/local/bin/defensia-agent && chmod 755 /usr/local/bin/defensia-agent && echo "defensia-agent: binary restored from backup" >&2; fi'
ExecStart=/usr/local/bin/defensia-agent start
Restart=always
RestartSec=10
# Disable systemd's start-limit so rapid successive auto-updates never
# leave the service stuck in start-limit-hit state.
StartLimitIntervalSec=0
StandardOutput=journal
StandardError=journal
SyslogIdentifier=defensia-agent
Environment=DEFENSIA_CONFIG=/etc/defensia/config.json

# Needs root for iptables — hardening where possible
NoNewPrivileges=no
PrivateTmp=yes
ProtectSystem=false
ProtectHome=yes

[Install]
WantedBy=multi-user.target
